# ansible/roles/workspace/templates/workspace-git-sync.timer.j2
[Unit]
Description=Hourly workspace git sync ({{ item.agent_id }})

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
