# ansible/roles/workspace/templates/workspace-git-sync.service.j2
[Unit]
Description=Sync OpenClaw workspace to git

[Service]
Type=oneshot
ExecStart=%h/.local/bin/workspace-git-sync.sh
Environment=SSH_AUTH_SOCK=
