# ansible/roles/workspace/templates/workspace-git-sync.service.j2
[Unit]
Description=Sync OpenClaw workspace ({{ item.agent_id }}) to git

[Service]
Type=oneshot
ExecStart=%h/.local/bin/workspace-git-sync-{{ item.agent_id }}.sh
Environment=SSH_AUTH_SOCK=
