# ansible/roles/plugins/templates/mcp-auth-proxy.service.j2
[Unit]
Description=MCP Auth Proxy (credential injection for Codex, Claude Code, and Pi containers)
After=network.target
Before=openclaw-gateway.service

[Service]
Type=simple
ExecStart=/usr/bin/node /home/ubuntu/.openclaw/mcp-auth-proxy.js
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
Environment=CODEX_PROXY_LISTEN={{ codex_proxy_gateway_ip }}
Environment=CODEX_PROXY_PORT={{ openclaw_mcp_adapter.codex_proxy_port }}
Environment=CODEX_AUTH_PATH=/home/ubuntu/.codex/auth.json
Environment=ANTHROPIC_AUTH_PATH=/home/ubuntu/.openclaw/anthropic-auth-token
Environment=GITHUB_TOKENS_DIR=/home/ubuntu/.openclaw/github-tokens

[Install]
WantedBy=default.target
